# deploy/immoscan-ventes.service
[Unit]
Description=ImmoScan - Pipeline ventes (scraping + scoring + alertes Telegram)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ImmoScan
Environment=PATH=/opt/ImmoScan/venv/bin:/usr/bin:/bin
ExecStart=/opt/ImmoScan/venv/bin/python scripts/run_pipeline_live.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300

